# com.research-pipeline.prefetch.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.research-pipeline.prefetch</string>
    <!-- ✅した論文のOA PDFを日中4回、自動でpapers/に先取り（Slackタップも同期）。
         深掘り取り込み(LLM)はしない＝非LLM・非課金のPDF確保だけ。
         有効化: launchctl load ~/Library/LaunchAgents/com.research-pipeline.prefetch.plist -->
    <key>ProgramArguments</key>
    <array>
        <string>/Users/tonn/matlab_pyenv312/bin/python3</string>
        <string>/Users/tonn/Desktop/Project/research-pipeline/promote_check.py</string>
        <string>--prepare</string>
    </array>
    <key>WorkingDirectory</key>
    <string>/Users/tonn/Desktop/Project/research-pipeline</string>
    <key>StartCalendarInterval</key>
    <array>
        <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>30</integer></dict>
        <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>30</integer></dict>
    </array>
    <key>RunAtLoad</key>
    <false/>
    <key>StandardOutPath</key>
    <string>/Users/tonn/Desktop/Project/research-pipeline/prefetch.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/tonn/Desktop/Project/research-pipeline/prefetch.log</string>
</dict>
</plist>
